# etc/cron.d/sandbox
PATH=/usr/sbin:/usr/bin:/sbin:/bin
0 * * * * root sandbox-upgrade >/dev/null 2>/dev/null
